# deploy/crawler/systemd/qbu-crawler-daily-submit.service
[Unit]
Description=Qbu-Crawler daily submit workflow
Wants=network-online.target
After=network-online.target qbu-crawler.service

[Service]
Type=oneshot
WorkingDirectory=/opt/qbu-crawler
EnvironmentFile=/etc/qbu-crawler.env
ExecStart=/usr/bin/env uv run python main.py workflow daily-submit

[Install]
WantedBy=multi-user.target
