# deploy/crawler/systemd/qbu-crawler-daily-submit.timer
[Unit]
Description=Run Qbu-Crawler daily submit workflow every day

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=qbu-crawler-daily-submit.service

[Install]
WantedBy=timers.target
